# packaging/build_files/etc/supervisord.d/02-EndToEndWorker.conf
[program:end-to-end-tester]
command=/usr/bin/python3 /usr/local/sbin/endtester.py

[program:db-recorder]
command=/usr/bin/python3 /usr/local/sbin/dbrecorder.py
