# packaging/build_files/etc/supervisord.d/03-mariadbstart.conf
[program:dbstart]
command=bash -c "/usr/local/sbin/dbstart.sh"
